# backend/fitmeals.service
[Unit]
Description=FitMeals Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/fitmeals-backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
